# python/gemini/ai.service
[Unit]
Description=AI Web Application
After=network.target

[Service]
# 使用低权限的专用用户
User=www-data
# 使用低权限的专用组
Group=www-data              

# 设置工作目录
WorkingDirectory=/www

# 限制最大虚拟内存为 1 GB
ExecStartPre=/bin/bash -c "ulimit -v 1048576"  
# 限制最大文件描述符数
ExecStartPre=/bin/bash -c "ulimit -n 65536"    

# 使用 taskset 设置 CPU 亲和性，绑定到 CPU 核心 0 和 1
ExecStart=/bin/taskset -c 0,1 /www/venv/bin/hypercorn -w 4 -b 0.0.0.0:5000 app:app  

# 配置环境变量，确保使用虚拟环境的路径
Environment=PATH=/www/venv/bin    
# 配置 Python 输出不被缓冲
Environment=PYTHONUNBUFFERED=1   

# 服务失败时重启
Restart=on-failure                
# 重启间隔设置为 5 秒
RestartSec=5                      
# 重启频率不做限制
StartLimitIntervalSec=0           

# 使用私有的临时文件系统
PrivateTmp=true                   
# 禁止子进程获得额外权限
NoNewPrivileges=true              
# 限制对系统文件的修改权限
ProtectSystem=full                
# 限制对用户目录的访问
ProtectHome=yes                   
# 限制网络协议类型为 IPv4
RestrictAddressFamilies=AF_INET   
# 仅授予绑定端口权限
AmbientCapabilities=CAP_NET_BIND_SERVICE  

# 限制最大文件描述符数
LimitNOFILE=65536                 
# 限制最大进程数
LimitNPROC=1000                   
# 限制最大内存使用为 1 GB
LimitAS=1G                        

[Install]
# 设置为多用户目标，开机启动
WantedBy=multi-user.target        
