# xray/xray.service
# 服务单元描述
[Unit]
Description=Xray Proxy Service
After=network-online.target nss-lookup.target
Wants=network-online.target

# 服务行为配置
[Service]
Type=simple
User=xvpn
Group=xvpn
WorkingDirectory=$INSTALL_DIR
ExecStart=$INSTALL_DIR/xray run -config $INSTALL_DIR/config.json
Restart=on-failure
RestartSec=3

# 设置最大文件描述符数量，提高并发连接处理能力
LimitNOFILE=1048576
# 设置最大进程数，限制资源使用，防止 fork 滥用
LimitNPROC=4096
# 禁止生成核心转储文件，提升安全性
LimitCORE=0
# 设置线程栈大小为 8MB，防止内存耗尽
LimitSTACK=8M

# 禁止服务提升权限（如 setuid）
NoNewPrivileges=true
# 使用独立的 /tmp 空间，防止与其他进程共享
PrivateTmp=true
# 禁止访问主机设备（如 /dev/sda）
PrivateDevices=true
# 系统目录（如 /usr /etc）设为只读，增加保护
ProtectSystem=strict
# 禁止访问 /home、/root 等主目录
ProtectHome=true
# 禁止加载或卸载内核模块
ProtectKernelModules=true
# 禁止访问 cgroup 控制组
ProtectControlGroups=true
# 禁止修改系统时间
ProtectClock=true
# 禁止访问 /proc/kmsg（内核日志）
ProtectKernelLogs=true
# 隐藏除自身外的 /proc 中其他进程信息
ProtectProc=invisible
# 限制 /proc 的子系统访问为 pid（增强隔离）
ProcSubset=pid

# 显式允许写入这些目录（其余系统只读）
ReadWritePaths=$INSTALL_DIR

# 仅保留必要的能力：绑定端口和操作 TUN 接口
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_ADMIN
# 提供运行所需的能力集，支持透明代理或 TUN
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_ADMIN

# 屏蔽敏感系统调用：挂载、重启、加载模块等
SystemCallFilter=~@mount @reboot @swap @module
# 仅允许当前 CPU 架构的系统调用，提升隔离
SystemCallArchitectures=native
# 禁止访问这些路径
InaccessiblePaths=/bin /sbin /boot /home /media /srv /root /var /opt /sys /usr /www /proc /lost+found /lib /lib64

# 额外的内存保护，防止恶意代码注入
MemoryDenyWriteExecute=true
# 限制内存锁定，防止进程锁定过多内存，避免内存消耗过多
LimitMEMLOCK=0

# 设置为默认多用户系统目标的一部分（随系统启动）
[Install]
WantedBy=multi-user.target
